# cron.d/osg-update-certs
# Cron job running by default every 1 hour.
# The lock file can be enabled or disabled via a
# service osg-update-certs-cron start
# chkconfig osg-update-certs-cron on

# Note the lock file not existing is success hence the the slightly odd logic
# below.

# The random sleep time is up to 45 minutes (2700 seconds).  The reason is so that
# not every single site gets their certificates at the exact same time and hits
# the CA cert server really hard.

10 * * * *    root    [ ! -f /var/lock/subsys/osg-update-certs-cron ] || /usr/sbin/osg-update-certs --random-sleep 2700 --called-from-cron >/dev/null 2>&1
